# track_it.service
[Unit]
Description=A Daemon for File Tracking
Wants= 


[Service]
ExecStart=/usr/local/bin/track_it
